# vv.conf
[program:vv]
command=python3 /home/curator2/vv/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/vv/test.err.log
stdout_logfile=/home/curator2/vv/test.out.log